# docker/supervisor.cron.conf
[program:cron]
command=/usr/sbin/crond -f -L 15
autostart=true
autorestart=true
